# centos7-phpenv-and-httpd/supervisord-fpm7.3.0-service.ini
[program:php-fpm7.3.0]
command=systemctl start php-fpm7.3.0