# shells/templates/supervisord-fpm-service.ini.tpl
[program:php-fpm${SERVICE_NAME}]
command=systemctl start php-fpm${SERVICE_NAME}