# shells/templates/supervisord-fpm-service.ini.tpl
[program:php-fpm${SERVICE_VERSION}]
command=systemctl start php-fpm${SERVICE_VERSION}