# web/supervisord-fpm7.1.20-service.ini
[program:php-fpm7.1.20]
command=systemctl start php-fpm7.1.20